# roles/confluence/templates/confluence.service.j2
# {{ ansible_managed }}

[Unit]
Description=Atlassian Confluence
After=syslog.target network.target

[Service]
Type=forking
User={{ confluence_owner }}
{% if confluence_service_environment is defined %}
{%    for cse in confluence_service_environment %}
Environment={{cse}}
{%    endfor %}
{% endif %}
PIDFile={{ confluence_install_dir }}/work/catalina.pid
ExecStart={{ confluence_install_dir }}/bin/start-confluence.sh
ExecStop={{ confluence_install_dir }}/bin/stop-confluence.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
